# deployment/chatbot-backend.service
[Unit]
Description=Chatbot Backend (FastAPI)
After=network.target qdrant.service
Requires=network.target

[Service]
Type=simple
User=chatbot
Group=chatbot
WorkingDirectory=/opt/chatbot/backend
Environment="PATH=/opt/chatbot/backend/venv/bin"
EnvironmentFile=/opt/chatbot/backend/.env
ExecStart=/opt/chatbot/backend/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chatbot-backend

# 리소스 제한 (선택사항)
LimitNOFILE=65535
LimitNPROC=4096

[Install]
WantedBy=multi-user.target

